# pull_deploy/linux/build/Docker/crontab
* *  *   *   *     $HOME/update-service.sh >> /var/log/cron.log 2>&1 && $HOME/stop-service.sh >> /var/log/cron.log 2>&1
