# pull_deploy/linux/Docker/crontab
* *  *   *   *     /home/$BASE_IMAGE/update-service.sh >> /var/log/cron.log 2>&1
